# files/etc/profile.d/alias.sh
alias cget="wget -q -O -"
alias portal="wget -q -O - http://detectportal.firefox.com/success.txt"
alias lh='ls -alhF --color=auto'
